# supervisor/TelegramBot-Queue-TelegramLimitedApiRequest.ini
[program:TelegramBot-Queue-TelegramLimitedApiRequest]
user = www
process_name = TelegramLimitedApiRequest_%(process_num)02d
numprocs = 4
directory = /www/wwwroot/TelegramBot
command = /www/server/php/bin/php artisan queue:work --queue=TelegramLimitedApiRequest --tries=3 --max-time=3600
autostart = true
autorestart = true
stopasgroup = true
killasgroup = true
stopwaitsecs = 3600
redirect_stderr = true
stdout_logfile = /www/wwwroot/TelegramBot/storage/logs/TelegramLimitedApiRequest.queue.%(process_num)02d.out.log
stderr_logfile = /www/wwwroot/TelegramBot/storage/logs/TelegramLimitedApiRequest.queue.%(process_num)02d.err.log
