# supervisor/TelegramBot-Queue-default.ini
[program:TelegramBot-Queue-default]
user=www
process_name=%(program_name)s_%(process_num)02d
directory=/www/wwwroot/TelegramBot
command=/www/server/php/bin/php artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=4
redirect_stderr=true
stdout_logfile=/www/wwwroot/TelegramBot/storage/logs/default.queue.out.log
stderr_logfile=/www/wwwroot/TelegramBot/storage/logs/default.queue.err.log
